# docker/legacy/cron/crontab
# Cron jobs for legacy game server
# Turn processing runs every minute

# Game turn processing (every minute)
* * * * * curl -sf http://legacy-nginx/api.php?action=turn > /dev/null 2>&1

# Cleanup old sessions (daily at 3 AM)
0 3 * * * find /tmp/sessions -type f -mtime +1 -delete > /dev/null 2>&1

# Clear old logs (weekly on Sunday at 4 AM)
0 4 * * 0 find /var/www/html/d_log -type f -mtime +30 -delete > /dev/null 2>&1
